# meshtastic-gps-forwarder.service
[Unit]
Description=Meshtastic GPS Position Injector
After=network.target gpsd.service
Wants=gpsd.service

[Service]
Type=simple
User=jeff
Group=dialout
WorkingDirectory=/opt/rpi-clock
ExecStart=/opt/rpi-clock/venv/bin/python3 -u /opt/rpi-clock/gps-meshtastic-inject.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=meshtastic-gps

[Install]
WantedBy=multi-user.target
